# final.target
[Unit]
Description=Last Steps
Wants=dwm.target
AllowIsolate=true

[Install]
Alias=default.target

