# mailer.service
[Unit]
Description=Ip mailer daemon

[Service]
ExecStart=/home/bighops/bin/mailer
RemainAfterExit=yes

[Install]
WantedBy=mystuff.target

