# mystuff.target
[Unit]
Description=Stuff to run in DWM.
After=xorg.target
