# taplist.service
[Unit]
Description=Refresh Taplist

[Service]
Environment=DISPLAY=:0
ExecStart=/usr/bin/dwb -x 'reload;;'
Restart=on-failure
