# xset.timer
[Unit]
Description=Make sure xset gets set

[Timer]
OnUnitInactiveSec=10m
OnStartupSec=1m

[Install]
WantedBy=mystuff.target
